# lib/hosting/assets/systemd.service
[Unit]
Description=Hosting account {{name}}

[Service]
Type=simple
ExecStart={{command}}

[Install]
WantedBy=multi-user.target
